# ergod.service
[Unit]
Description=computer ergonomics daemon

[Service]
Type=simple
ExecStart=.../ergod
