# ergod.service
[Unit]
Description=computer ergonomics daemon

[Service]
Type=simple
ExecStart=!/usr/local/sbin/ergod

[Install]
WantedBy=default.target
